# trojan-go-s6/services.d/trojan-go/run
#!/bin/sh

exec 2>&1
exec /usr/local/bin/trojan-go -config /etc/trojan-go/config.json