# puma/runit/log/run
#!/bin/bash

exec chpst -u app svlogd -tt $APP_CURRENT_PATH/log/puma
